# systemd/murmur-daemon.service
[Unit]
Description=Murmur Voice Transcription Daemon
Documentation=https://github.com/abrahamvdl/murmur
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/murmur-daemon
Restart=on-failure
RestartSec=5

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="QT_QPA_PLATFORM=wayland"
# Uncomment and adjust for your AMD GPU if needed
# Environment="HSA_OVERRIDE_GFX_VERSION=10.3.0"

# Resource limits
MemoryMax=4G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=murmur-daemon

[Install]
WantedBy=default.target
